# .config/systemd/user/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (v2026.2.1)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart="/usr/bin/node" "/home/kaar/.npm-global/lib/node_modules/openclaw/dist/index.js" gateway --port 18789
Restart=always
RestartSec=5
KillMode=process
Environment=HOME=/home/kaar
Environment="PATH=/home/kaar/.local/share/pnpm:/home/kaar/.config/nvm/current/bin:/home/kaar/.local/bin:/home/kaar/.npm-global/bin:/home/kaar/bin:/home/kaar/.nvm/current/bin:/home/kaar/.fnm/current/bin:/home/kaar/.volta/bin:/home/kaar/.asdf/shims:/home/kaar/.bun/bin:/usr/local/bin:/usr/bin:/bin"
Environment=OPENCLAW_GATEWAY_PORT=18789
Environment=OPENCLAW_GATEWAY_TOKEN=173eac79f339c5a8c0e57463fd7d6cfb4fab477c1641ea41
Environment="OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service"
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION=2026.2.1

[Install]
WantedBy=default.target
